# daemontools/log/run
#!/bin/sh
export SERVICE=$(basename $(dirname $PWD))
exec logger -t $SERVICE
